# xbmc-es-switch.service
[Unit]
Description = Starts the XBMC-ES-Switch program
After = remote-fs.target

[Service]
Type = simple
ExecStart = /usr/bin/xbmc-es-switch

[Install]
WantedBy = multi-user.target
